# crontab
* * * * * ~/src/capture/capture.sh > /dev/null
